# os/linux/gnome/autostart/DisableTracker.desktop
[Desktop Entry]
Type=Application
Name=Disable Tracker-Store
Comment=Disable Tacker-Store
Exec=tracker-control -r
X-GNOME_Autostart-enabled=true
